# rco-idv-client/client/scripts/usb.rules
KERNEL=="sd*[0-9]",SUBSYSTEMS=="usb", RUN+="/usr/local/bin/usb_mnt %k %n"
